# infrastructure/terraform/environments/development/files/config-sync.service
[Unit]
Description=Config Sync from DigitalOcean Spaces (Zero-SSH Pattern)
After=network-online.target docker.service
Wants=network-online.target
PartOf=config-sync.timer

[Service]
Type=oneshot
User=root
Group=root

# Environment Configuration
Environment="SPACES_BUCKET=ectropy-development-configs"
Environment="SPACES_REGION=sfo3"
Environment="CONFIG_DIR=/opt/ectropy"
Environment="LOG_FILE=/var/log/config-sync.log"
Environment="DOCKER_CONFIG=/opt/ectropy/.docker"

# AWS Credentials (provided by cloud-init)
EnvironmentFile=-/etc/default/config-sync

# Working Directory
WorkingDirectory=/opt/ectropy

# Execute sync script
ExecStart=/opt/ectropy/scripts/config-sync.sh

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=config-sync

# Security Settings
ProtectSystem=full
ProtectHome=yes
NoNewPrivileges=yes
PrivateTmp=yes

# Resource Limits
TimeoutStartSec=120s
Restart=on-failure
RestartSec=60s

[Install]
WantedBy=multi-user.target
