# infrastructure/terraform/environments/staging/files/config-sync.service
[Unit]
Description=Config Sync (fast path — pull from Spaces, flag restart if changed)
Documentation=https://github.com/luhtech/Ectropy/blob/main/.roadmap/PHASE_P3_ZERO_SSH_TERRAFORM_STRATEGY_2026-02-14.json
After=network-online.target docker.service
Wants=network-online.target
PartOf=config-sync.timer

[Service]
Type=oneshot
User=root
Group=root

# RuntimeDirectory: tmpfs /run/config-sync/ (flag-file handoff to config-restart.service)
# Created by systemd before ExecStart; owned by root:root, mode 0755
RuntimeDirectory=config-sync
RuntimeDirectoryPreserve=yes
RuntimeDirectoryMode=0755

# Environment Configuration
Environment="SPACES_BUCKET=ectropy-staging-configs"
Environment="SPACES_REGION=sfo3"
Environment="CONFIG_DIR=/opt/ectropy"
Environment="LOG_FILE=/var/log/config-sync.log"
Environment="DOCKER_CONFIG=/opt/ectropy/.docker"

# AWS Credentials (provided by cloud-init or manually configured)
# Required: AWS_ACCESS_KEY_ID and AWS_SECRET_ACCESS_KEY
# DO NOT hardcode credentials - use environment variables or instance metadata
EnvironmentFile=-/etc/default/config-sync

# Working Directory
WorkingDirectory=/opt/ectropy

# Execute sync script in SYNC-ONLY mode (fast path, flag-based handoff)
# Restart phase runs in config-restart.service (triggered by config-restart.path)
ExecStart=/opt/ectropy/scripts/config-sync.sh --sync-only

# Logging Configuration
StandardOutput=journal
StandardError=journal
SyslogIdentifier=config-sync

# Security Settings (Fortune 500 Compliance)
# ProtectSystem=strict would prevent writing to /opt/ectropy
ProtectSystem=full
ProtectHome=yes
NoNewPrivileges=yes
PrivateTmp=yes

# Resource Limits — 60s budget for sync only (no docker compose up)
TimeoutStartSec=60s
Restart=on-failure
RestartSec=60s

[Install]
WantedBy=multi-user.target
